# root/app/certbot/supervisord.conf
[program:certbot-renew]
command = /app/certbot/renew.sh
autorestart = true
pidfile = /run/certbot-renew.pid
stdout_logfile = /var/log/certbot-renew.log
logfile_maxbytes = 50MB
priority = 2
